# lib/actor/fleet-heartbeat.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!-- kototama fleet heartbeat — durable schedule for continuous organism self-update.
     Install:  cp this → ~/Library/LaunchAgents/com.kototama.fleet-heartbeat.plist
               (edit ABSOLUTE paths + interval below), then
               launchctl load ~/Library/LaunchAgents/com.kototama.fleet-heartbeat.plist
     Off-switch (種をまく leash):  launchctl unload … ; the organisms fall back to their local logs.
     Each firing runs one --live pass: every actor's beat (regen) + ipfs name publish under its own key. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>            <string>com.kototama.fleet-heartbeat</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/bb</string>
    <string>/Users/junkawasaki/github/com-junkawasaki/orgs/com-junkawasaki/kototama/lib/actor/fleet-heartbeat.bb</string>
    <string>--once</string>
    <string>--live</string>
    <string>--root</string>
    <string>/Users/junkawasaki/github/com-junkawasaki/orgs/etzhayyim</string>
  </array>
  <!-- every 6 hours; raise/lower to taste. Each pass = ~50 ipfs name publishes. -->
  <key>StartInterval</key>    <integer>21600</integer>
  <key>RunAtLoad</key>        <true/>
  <key>StandardOutPath</key>  <string>/tmp/kototama-fleet-heartbeat.log</string>
  <key>StandardErrorPath</key><string>/tmp/kototama-fleet-heartbeat.err</string>
</dict>
</plist>
